# ufw-blocklist-update.service
[Unit]
Description=Pull online blocklists and refresh UFW ipset
After=network-online.target ufw.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/update-blocklist-sources.sh
